# deploy/strom-agent.service
[Unit]
Description=Strom node agent
After=network-online.target nut-server.service strom-firstboot.service
Wants=network-online.target
RequiresMountsFor=/var/lib/strom

[Service]
Type=simple
ExecStart=/usr/local/bin/strom-agent --config-dir /etc/nut --listen :80
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
